# rebuild_snapshot_2026-04-21_031038/openclaw-backup.service
[Unit]
Description=Backup OpenClaw workspace to GitHub when changed

[Service]
Type=oneshot
ExecStart=/root/.openclaw/scripts/backup_openclaw_to_github.sh
