# rebuild_snapshot_2026-04-15_031038/polymarket-direct-monitor.service
[Unit]
Description=Direct Polymarket monitor to Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/.openclaw/scripts/polymarket_direct_monitor.sh
WorkingDirectory=/root/.openclaw/workspace
